# backend/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:cron]
command=/usr/sbin/cron -f
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fastapi]
directory=/app/src/
# %(ENV_NGX_PORT)s = supervisord substitution from its own env (NOT bash
# $NGX_PORT — bash sub doesn't fire under supervisord). Falls back to
# 13491 if NGX_PORT isn't set so a bare container still boots.
environment=NGX_PORT="%(ENV_NGX_PORT)s"
command=/app/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port %(ENV_NGX_PORT)s --workers 4
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
